# setup/etc/cron.d/issabel-portknock.cron
0 * * * * root /usr/bin/issabel-portknock-cleanup