# setup/etc/cron.daily/renewssl
#!/bin/bash
if [ -f /etc/letsencrypt/values ];
then
    /usr/share/issabel/privileged/ssl_certbot renewcertificate >/dev/null
fi
