# setup/etc/cron.daily/renewssl
#!/bin/bash
if [ -f /etc/letsencrypt/values ];
then
/usr/bin/certbot renew  --quiet
fi
